# conf/supervisor/supervisor.d/c2c.conf
[program:c2c]
command=/usr/local/sbin/c2c.py
autostart=true
autorestart=true
