# conf/supervisor/c2c.d.conf
[program:c2c.d]
command=/usr/local/sbin/c2c.d.py
autostart=true
autorestart=true
user=c2c
